# amxa-desktop-extra-1.0/etc/profile.d/hide-home-icon.sh
if test "$(cat /etc/puavo/hosttype)" = "laptop"; then
  gsettings set org.gnome.nautilus.desktop home-icon-visible false
else
  gsettings set org.gnome.nautilus.desktop home-icon-visible true
fi

